# packages/operbot/operbot-165-py3-none-any.whl/operbot-165.data/data/operbot/operbot.service
[Unit]
Description=A contribution back to society
After=multi-user.target

[Service]
#Type=simple
User=operbot
Group=operbot
DynamicUser=yes
StateDirectory=operbot
LogsDirectory=operbot
CacheDirectory=operbot
ExecStart=/usr/local/bin/operbotd
CapabilityBoundingSet=CAP_NET_RAW

[Install]
WantedBy=multi-user.target
