# packages/gcid/gcid-160-py3-none-any.whl/gcid-160.data/data/gcid/gcid.service
[Unit]
Description=medicine care, poison genocide
After=multi-user.target

[Service]
#Type=simple
User=gcid
Group=gcid
DynamicUser=yes
StateDirectory=gcid
LogsDirectory=gcid
CacheDirectory=gcid
ExecStart=/usr/local/bin/gcided
CapabilityBoundingSet=CAP_NET_RAW

[Install]
WantedBy=multi-user.target
